# files/systemd/system/setup-dialout-group.service
[Unit]
Description=Set up dialout group and add all users to it
ConditionPathExists=!/etc/setup-dialout-group
After=multi-user.target

[Service]
Type=oneshot
# INFO: dialout group exists by default and has to be copied to /etc/group in bootc systems from /usr/lib/group
ExecStart=/usr/bin/bash -c "grep -E '^dialout:' /usr/lib/group >> /etc/group; for user in $(getent passwd | grep 'home' | awk -F: '{print $1}'); do usermod -aG dialout $user; done"
ExecStart=/usr/bin/touch /etc/setup-dialout-group
ExecStart=/usr/bin/udevadm control --reload-rules

[Install]
WantedBy=multi-user.target
